# spec/fixtures/unit/puppetx/parsers/crontab/misplaced_puppet_name
# Puppet Name: first
ENV_VAR=VAR_VAL
# Puppet Name: second
* * * * * /bin/true
