# spec/fixtures/unit/puppetx/parsers/crontab/repeated_puppet_name
# Puppet Name: first
# Puppet Name: second
* * * * * /bin/true
